# src/mcp/wheelhive-api.service
[Unit]
Description=WheelHive API Server
After=network.target

[Service]
Type=simple
User=steve
Group=steve
WorkingDirectory=/home/steve/code/wheelhive
Environment="PATH=/home/steve/code/wheelhive/.venv/bin"
ExecStart=/home/steve/code/wheelhive/.venv/bin/python /home/steve/code/wheelhive/src/mcp/wheelhive_api.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wheelhive-api

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/steve/code/wheelhive

[Install]
WantedBy=multi-user.target
